# opal2-vm/salt/tmpwatch/tmpwatch.cron
0 3 * * *  root  tmpwatch 336 /vol_c/eemtws-jobs/ >/dev/null 2>&1

